# deploy/obsidian-gateway.service
# systemd --user unit for the @stable uv-tool deploy (daily auto-update: deploy/obsidian-gateway-update.*).
# Install: uv tool install --from git+https://github.com/fszalaj/obsidian-gateway@stable obsidian-gateway
#   cp {vaults,tokens}.example.yaml ~/obsidian-gateway/ (edit; chmod 0600 tokens.yaml)
#   cp deploy/obsidian-gateway.service ~/.config/systemd/user/ && systemctl --user enable --now obsidian-gateway
#   loginctl enable-linger $USER   # Tailscale Serve fronts the loopback bind.

[Unit]
Description=Obsidian MCP gateway (filesystem/git-native)
After=network-online.target

[Service]
Type=simple
Environment=OBSIDIAN_GATEWAY_HOST=127.0.0.1
Environment=OBSIDIAN_GATEWAY_PORT=8765
Environment=OBSIDIAN_GATEWAY_PATH=/mcp/
Environment=OBSIDIAN_GATEWAY_VAULTS=%h/obsidian-gateway/vaults.yaml
Environment=OBSIDIAN_GATEWAY_TOKENS=%h/obsidian-gateway/tokens.yaml
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=%h/.local/bin/obsidian-gateway
Restart=on-failure
RestartSec=3

NoNewPrivileges=yes
CapabilityBoundingSet=
AmbientCapabilities=
RestrictSUIDSGID=yes
RestrictRealtime=yes
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
ProtectClock=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK
UMask=0007

# Stronger isolation, but needs unprivileged userns (can fail to start without it).
# Probe per host: systemd-run --user -p <Directive> /bin/true -> rc 0 = supported.
# Then uncomment what passes and add a ReadWritePaths per vault repo_root.
# PrivateTmp=yes
# PrivateDevices=yes
# ProtectKernelTunables=yes
# ProtectKernelModules=yes
# ProtectKernelLogs=yes
# ProtectControlGroups=yes
# ProtectHostname=yes
# ProtectProc=invisible
# RemoveIPC=yes
# ProtectSystem=strict
# ReadWritePaths=%h/obsidian-vaults/your-repo

[Install]
WantedBy=default.target
